# root/etc/s6-overlay/s6-rc.d/svc-tachidesk/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mv /build/Suwayomi-Server.jar /app/suwayomi-server/Suwayomi-Server.jar 
# export suwayomi-server_DBPATH=/config

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 8085" \
    s6-setuidgid abc java -Dsuwayomi.suwayomi-server.config.server.rootDir="/config" -jar /app/suwayomi-server/Suwayomi-Server.jar \
    -nobrowser -data=/config
